# debian/postinst
#!/bin/bash
set -e

echo "[AutoRip] Configuration du service AutoRip..."

# Activation des règles udev
udevadm control --reload-rules
udevadm trigger

# Activation du démon
systemctl daemon-reload
systemctl enable autorip.service || true

# Création du dossier temporaire
mkdir -p /tmp/autorip

echo "[AutoRip] Configuration terminée !"
exit 0
