# debian/postinst
#!/bin/bash
set -e

echo "[AutoRip] Installation des dépendances..."

# Ajout du dépôt MakeMKV si nécessaire
if ! command -v makemkvcon &> /dev/null; then
    echo "[AutoRip] Installation de MakeMKV..."
    apt-get install -y curl gnupg ca-certificates python3
    KEY_FP=$(curl -fsSL "https://api.launchpad.net/1.0/~heyarje/+archive/makemkv-beta" \
        | python3 -c "import sys,json; print(json.load(sys.stdin)['signing_key_fingerprint'])")
    gpg --keyserver hkp://keyserver.ubuntu.com:80 --recv-keys "$KEY_FP"
    gpg --export "$KEY_FP" > /etc/apt/trusted.gpg.d/makemkv-beta.gpg
    echo "deb https://ppa.launchpadcontent.net/heyarje/makemkv-beta/ubuntu jammy main" \
        > /etc/apt/sources.list.d/makemkv-beta.list
    apt-get update
    apt-get install -y makemkv-bin makemkv-oss
fi

# Installation libdvdcss2
if ! dpkg -l | grep -q libdvdcss2; then
    echo "[AutoRip] Installation de libdvdcss2..."
    apt-get install -y libdvdcss2 || true
fi

# Activation des règles udev
udevadm control --reload-rules
udevadm trigger

# Activation du démon
systemctl daemon-reload
systemctl enable autorip.service || true

# Création du dossier temporaire
mkdir -p /tmp/autorip

echo "[AutoRip] Installation terminée !"
exit 0
