# etc/udev/rules.d/99-dvd-autorip.rules
# Règle udev : détection insertion DVD et lancement du rip automatique
ACTION=="change", KERNEL=="sr[0-9]*", ENV{ID_CDROM_MEDIA_DVD}=="1", RUN+="/usr/bin/python3 /srv/autorip/autorip_daemon.py --device /dev/%k"
